# cron.daily/toggl_integrate
#!/bin/sh

EXE=/opt/odoo_scripts/scripts/toggl_integrate.py
LOGFILE=/var/log/toggl_integrate.log

echo "=============== [$(date)] ===============" >> $LOGFILE
$EXE -u admin -p password >> $LOGFILE

